# todoist-sync/todoist-sync.timer
[Unit]
Description=Timer for todoist sync on Startup

[Timer]
OnBootSec=1min
Unit=todoist-sync.service

[Install]
WantedBy=multi-user.target